# scripts/deployment-status.timer
[Unit]
Description=Timer to schedule deployment-status.service to run

[Timer]
OnCalendar=minutely
